# docker/backup/crontab
# Daily backup at 03:00 UTC
0 3 * * * /usr/local/bin/backup.sh
